# rootfs/var/lib/dpkg/info/dovecot-core.conffiles
/etc/default/dovecot
/etc/dovecot/README
/etc/init.d/dovecot
/etc/pam.d/dovecot
