# rootfs/var/lib/dpkg/info/fail2ban.list
/.
/usr
/usr/share
/usr/share/fail2ban
/usr/share/fail2ban/common
/usr/share/fail2ban/common/__init__.py
/usr/share/fail2ban/common/exceptions.py
/usr/share/fail2ban/common/helpers.py
/usr/share/fail2ban/common/protocol.py
/usr/share/fail2ban/common/version.py
/usr/share/fail2ban/client
/usr/share/fail2ban/client/__init__.py
/usr/share/fail2ban/client/actionreader.py
/usr/share/fail2ban/client/beautifier.py
/usr/share/fail2ban/client/configparserinc.py
/usr/share/fail2ban/client/configreader.py
/usr/share/fail2ban/client/configurator.py
/usr/share/fail2ban/client/csocket.py
/usr/share/fail2ban/client/fail2banreader.py
/usr/share/fail2ban/client/filterreader.py
/usr/share/fail2ban/client/jailreader.py
/usr/share/fail2ban/client/jailsreader.py
/usr/share/fail2ban/server
/usr/share/fail2ban/server/__init__.py
/usr/share/fail2ban/server/action.py
/usr/share/fail2ban/server/actions.py
/usr/share/fail2ban/server/asyncserver.py
/usr/share/fail2ban/server/banmanager.py
/usr/share/fail2ban/server/datedetector.py
/usr/share/fail2ban/server/datetemplate.py
/usr/share/fail2ban/server/faildata.py
/usr/share/fail2ban/server/failmanager.py
/usr/share/fail2ban/server/failregex.py
/usr/share/fail2ban/server/filter.py
/usr/share/fail2ban/server/filtergamin.py
/usr/share/fail2ban/server/filterpoll.py
/usr/share/fail2ban/server/filterpyinotify.py
/usr/share/fail2ban/server/iso8601.py
/usr/share/fail2ban/server/jail.py
/usr/share/fail2ban/server/jails.py
/usr/share/fail2ban/server/jailthread.py
/usr/share/fail2ban/server/mytime.py
/usr/share/fail2ban/server/server.py
/usr/share/fail2ban/server/ticket.py
/usr/share/fail2ban/server/transmitter.py
/usr/share/fail2ban/testcases
/usr/share/fail2ban/testcases/__init__.py
/usr/share/fail2ban/testcases/actionstestcase.py
/usr/share/fail2ban/testcases/actiontestcase.py
/usr/share/fail2ban/testcases/banmanagertestcase.py
/usr/share/fail2ban/testcases/clientreadertestcase.py
/usr/share/fail2ban/testcases/datedetectortestcase.py
/usr/share/fail2ban/testcases/dummyjail.py
/usr/share/fail2ban/testcases/failmanagertestcase.py
/usr/share/fail2ban/testcases/filtertestcase.py
/usr/share/fail2ban/testcases/misctestcase.py
/usr/share/fail2ban/testcases/samplestestcase.py
/usr/share/fail2ban/testcases/servertestcase.py
/usr/share/fail2ban/testcases/sockettestcase.py
/usr/share/fail2ban/testcases/utils.py
/usr/share/fail2ban/fail2ban-0.8.13.egg-info
/usr/share/doc
/usr/share/doc/fail2ban
/usr/share/doc/fail2ban/README.md
/usr/share/doc/fail2ban/run-rootless.txt
/usr/share/doc/fail2ban/TODO
/usr/share/doc/fail2ban/TODO.Debian
/usr/share/doc/fail2ban/copyright
/usr/share/doc/fail2ban/examples
/usr/share/doc/fail2ban/examples/ipmasq-ZZZzzz_fail2ban.rul
/usr/share/doc/fail2ban/examples/jail.conf.gz
/usr/share/doc/fail2ban/DEVELOP.gz
/usr/share/doc/fail2ban/changelog.Debian.gz
/usr/share/doc/fail2ban/NEWS.Debian.gz
/usr/share/doc/fail2ban/changelog.gz
/usr/share/doc/fail2ban/README.Debian.gz
/usr/share/python
/usr/share/python/runtime.d
/usr/share/python/runtime.d/fail2ban.rtupdate
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/fail2ban-client.1.gz
/usr/share/man/man1/fail2ban-regex.1.gz
/usr/share/man/man1/fail2ban-server.1.gz
/usr/share/man/man1/fail2ban.1.gz
/usr/share/man/man1/jail.conf.10.gz
/usr/bin
/usr/bin/fail2ban-client
/usr/bin/fail2ban-server
/usr/bin/fail2ban-regex
/etc
/etc/fail2ban
/etc/fail2ban/fail2ban.conf
/etc/fail2ban/filter.d
/etc/fail2ban/filter.d/3proxy.conf
/etc/fail2ban/filter.d/apache-auth.conf
/etc/fail2ban/filter.d/apache-badbots.conf
/etc/fail2ban/filter.d/apache-common.conf
/etc/fail2ban/filter.d/apache-modsecurity.conf
/etc/fail2ban/filter.d/apache-nohome.conf
/etc/fail2ban/filter.d/apache-noscript.conf
/etc/fail2ban/filter.d/apache-overflows.conf
/etc/fail2ban/filter.d/assp.conf
/etc/fail2ban/filter.d/asterisk.conf
/etc/fail2ban/filter.d/common.conf
/etc/fail2ban/filter.d/courierlogin.conf
/etc/fail2ban/filter.d/couriersmtp.conf
/etc/fail2ban/filter.d/cyrus-imap.conf
/etc/fail2ban/filter.d/dovecot.conf
/etc/fail2ban/filter.d/dropbear.conf
/etc/fail2ban/filter.d/ejabberd-auth.conf
/etc/fail2ban/filter.d/exim-common.conf
/etc/fail2ban/filter.d/exim-spam.conf
/etc/fail2ban/filter.d/exim.conf
/etc/fail2ban/filter.d/freeswitch.conf
/etc/fail2ban/filter.d/groupoffice.conf
/etc/fail2ban/filter.d/gssftpd.conf
/etc/fail2ban/filter.d/horde.conf
/etc/fail2ban/filter.d/lighttpd-auth.conf
/etc/fail2ban/filter.d/mysqld-auth.conf
/etc/fail2ban/filter.d/nagios.conf
/etc/fail2ban/filter.d/named-refused.conf
/etc/fail2ban/filter.d/nginx-http-auth.conf
/etc/fail2ban/filter.d/nsd.conf
/etc/fail2ban/filter.d/openwebmail.conf
/etc/fail2ban/filter.d/pam-generic.conf
/etc/fail2ban/filter.d/perdition.conf
/etc/fail2ban/filter.d/php-url-fopen.conf
/etc/fail2ban/filter.d/postfix-sasl.conf
/etc/fail2ban/filter.d/postfix.conf
/etc/fail2ban/filter.d/proftpd.conf
/etc/fail2ban/filter.d/pure-ftpd.conf
/etc/fail2ban/filter.d/qmail.conf
/etc/fail2ban/filter.d/recidive.conf
/etc/fail2ban/filter.d/roundcube-auth.conf
/etc/fail2ban/filter.d/selinux-common.conf
/etc/fail2ban/filter.d/selinux-ssh.conf
/etc/fail2ban/filter.d/sendmail-auth.conf
/etc/fail2ban/filter.d/sendmail-reject.conf
/etc/fail2ban/filter.d/sieve.conf
/etc/fail2ban/filter.d/sogo-auth.conf
/etc/fail2ban/filter.d/solid-pop3d.conf
/etc/fail2ban/filter.d/squid.conf
/etc/fail2ban/filter.d/sshd-ddos.conf
/etc/fail2ban/filter.d/sshd.conf
/etc/fail2ban/filter.d/suhosin.conf
/etc/fail2ban/filter.d/uwimap-auth.conf
/etc/fail2ban/filter.d/vsftpd.conf
/etc/fail2ban/filter.d/webmin-auth.conf
/etc/fail2ban/filter.d/wuftpd.conf
/etc/fail2ban/filter.d/xinetd-fail.conf
/etc/fail2ban/action.d
/etc/fail2ban/action.d/apf.conf
/etc/fail2ban/action.d/badips.conf
/etc/fail2ban/action.d/blocklist_de.conf
/etc/fail2ban/action.d/bsd-ipfw.conf
/etc/fail2ban/action.d/complain.conf
/etc/fail2ban/action.d/dshield.conf
/etc/fail2ban/action.d/dummy.conf
/etc/fail2ban/action.d/firewallcmd-ipset.conf
/etc/fail2ban/action.d/firewallcmd-new.conf
/etc/fail2ban/action.d/hostsdeny.conf
/etc/fail2ban/action.d/ipfilter.conf
/etc/fail2ban/action.d/ipfw.conf
/etc/fail2ban/action.d/iptables-allports.conf
/etc/fail2ban/action.d/iptables-blocktype.conf
/etc/fail2ban/action.d/iptables-ipset-proto4.conf
/etc/fail2ban/action.d/iptables-ipset-proto6-allports.conf
/etc/fail2ban/action.d/iptables-ipset-proto6.conf
/etc/fail2ban/action.d/iptables-multiport-log.conf
/etc/fail2ban/action.d/iptables-multiport.conf
/etc/fail2ban/action.d/iptables-new.conf
/etc/fail2ban/action.d/iptables-xt_recent-echo.conf
/etc/fail2ban/action.d/iptables.conf
/etc/fail2ban/action.d/mail-buffered.conf
/etc/fail2ban/action.d/mail-whois-lines.conf
/etc/fail2ban/action.d/mail-whois.conf
/etc/fail2ban/action.d/mail.conf
/etc/fail2ban/action.d/mynetwatchman.conf
/etc/fail2ban/action.d/osx-afctl.conf
/etc/fail2ban/action.d/osx-ipfw.conf
/etc/fail2ban/action.d/pf.conf
/etc/fail2ban/action.d/route.conf
/etc/fail2ban/action.d/sendmail-buffered.conf
/etc/fail2ban/action.d/sendmail-common.conf
/etc/fail2ban/action.d/sendmail-whois-lines.conf
/etc/fail2ban/action.d/sendmail-whois.conf
/etc/fail2ban/action.d/sendmail.conf
/etc/fail2ban/action.d/shorewall.conf
/etc/fail2ban/action.d/ufw.conf
/etc/fail2ban/fail2ban.d
/etc/fail2ban/jail.d
/etc/fail2ban/jail.conf
/etc/bash_completion.d
/etc/bash_completion.d/fail2ban
/etc/logrotate.d
/etc/logrotate.d/fail2ban
/etc/default
/etc/default/fail2ban
/etc/init.d
/etc/init.d/fail2ban
/var
/var/run
