# rootfs/var/lib/dpkg/info/mysql-server-5.5.list
/.
/etc
/etc/logrotate.d
/etc/logrotate.d/mysql-server
/etc/logcheck
/etc/logcheck/ignore.d.paranoid
/etc/logcheck/ignore.d.paranoid/mysql-server-5_5
/etc/logcheck/ignore.d.server
/etc/logcheck/ignore.d.server/mysql-server-5_5
/etc/logcheck/ignore.d.workstation
/etc/logcheck/ignore.d.workstation/mysql-server-5_5
/etc/init.d
/etc/init.d/mysql
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.mysqld
/etc/mysql
/etc/mysql/conf.d
/etc/mysql/conf.d/mysqld_safe_syslog.cnf
/etc/mysql/debian-start
/usr
/usr/share
/usr/share/mysql
/usr/share/mysql/mysql_test_data_timezone.sql
/usr/share/mysql/mysqld_multi.server
/usr/share/mysql/errmsg-utf8.txt
/usr/share/mysql/echo_stderr
/usr/share/mysql/debian-start.inc.sh
/usr/share/mysql/ndb-config-2-node.ini
/usr/share/mysql/config.small.ini
/usr/share/mysql/config.medium.ini
/usr/share/mysql/config.huge.ini
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mysql_tzinfo_to_sql.1.gz
/usr/share/man/man1/innochecksum.1.gz
/usr/share/man/man1/resolve_stack_dump.1.gz
/usr/share/man/man1/resolveip.1.gz
/usr/share/man/man1/replace.1.gz
/usr/share/man/man1/perror.1.gz
/usr/share/man/man1/mysql_zap.1.gz
/usr/share/man/man1/mysqltest.1.gz
/usr/share/man/man1/mysql_setpermission.1.gz
/usr/share/man/man1/mysql_secure_installation.1.gz
/usr/share/man/man1/mysqlhotcopy.1.gz
/usr/share/man/man1/mysqld_safe.1.gz
/usr/share/man/man1/mysqld_multi.1.gz
/usr/share/man/man1/mysql_convert_table_format.1.gz
/usr/share/man/man1/mysqlbinlog.1.gz
/usr/share/man/man1/myisampack.1.gz
/usr/share/man/man1/myisamlog.1.gz
/usr/share/man/man1/myisamchk.1.gz
/usr/share/man/man1/msql2mysql.1.gz
/usr/share/doc
/usr/share/doc/mysql-server-5.5
/usr/share/doc/mysql-server-5.5/README.Debian.gz
/usr/share/doc/mysql-server-5.5/changelog.Debian.gz
/usr/share/doc/mysql-server-5.5/NEWS.Debian.gz
/usr/share/doc/mysql-server-5.5/changelog.gz
/usr/share/doc/mysql-server-5.5/copyright
/usr/share/doc/mysql-server-5.5/examples
/usr/share/doc/mysql-server-5.5/examples/my-medium.cnf.gz
/usr/share/doc/mysql-server-5.5/examples/my-large.cnf.gz
/usr/share/doc/mysql-server-5.5/examples/my-innodb-heavy-4G.cnf.gz
/usr/share/doc/mysql-server-5.5/examples/my-huge.cnf.gz
/usr/share/doc/mysql-server-5.5/examples/my-small.cnf
/usr/share/doc/mysql-server-5.5/mysqld.sym.gz
/usr/bin
/usr/bin/myisamchk
/usr/bin/myisamlog
/usr/bin/myisampack
/usr/bin/mysql_tzinfo_to_sql
/usr/bin/mysqlbinlog
/usr/bin/mysqltest
/usr/bin/perror
/usr/bin/replace
/usr/bin/resolve_stack_dump
/usr/bin/resolveip
/usr/bin/mysqlhotcopy
/usr/bin/mysqld_safe
/usr/bin/mysqld_multi
/usr/bin/mysql_zap
/usr/bin/mysql_setpermission
/usr/bin/mysql_secure_installation
/usr/bin/mysql_convert_table_format
/usr/bin/msql2mysql
/usr/lib
/usr/lib/mysql
/usr/lib/mysql/plugin
/usr/lib/mysql/plugin/adt_null.so
/usr/lib/mysql/plugin/auth.so
/usr/lib/mysql/plugin/auth_socket.so
/usr/lib/mysql/plugin/auth_test_plugin.so
/usr/lib/mysql/plugin/ha_example.so
/usr/lib/mysql/plugin/libdaemon_example.so
/usr/lib/mysql/plugin/mypluglib.so
/usr/lib/mysql/plugin/qa_auth_client.so
/usr/lib/mysql/plugin/qa_auth_interface.so
/usr/lib/mysql/plugin/qa_auth_server.so
/usr/lib/mysql/plugin/semisync_master.so
/usr/lib/mysql/plugin/semisync_slave.so
/usr/share/man/man1/mysqltest_embedded.1.gz
