# rootfs/var/lib/dpkg/info/nginx-common.list
/.
/etc
/etc/nginx
/etc/nginx/proxy_params
/etc/nginx/snippets
/etc/nginx/snippets/snakeoil.conf
/etc/nginx/snippets/fastcgi-php.conf
/etc/nginx/sites-available
/etc/nginx/sites-available/default
/etc/nginx/scgi_params
/etc/nginx/koi-win
/etc/nginx/nginx.conf
/etc/nginx/mime.types
/etc/nginx/koi-utf
/etc/nginx/uwsgi_params
/etc/nginx/sites-enabled
/etc/nginx/fastcgi.conf
/etc/nginx/fastcgi_params
/etc/nginx/conf.d
/etc/nginx/win-utf
/etc/logrotate.d
/etc/logrotate.d/nginx
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/nginx
/etc/default
/etc/default/nginx
/etc/init.d
/etc/init.d/nginx
/var
/var/log
/var/log/nginx
/var/www
/var/www/html
/var/lib
/var/lib/nginx
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/nginx.1.gz
/usr/share/nginx
/usr/share/nginx/html
/usr/share/nginx/html/index.html
/usr/share/doc
/usr/share/doc/nginx-common
/usr/share/doc/nginx-common/copyright
/usr/share/doc/nginx-common/NEWS.Debian.gz
/usr/share/doc/nginx-common/changelog.gz
/usr/share/doc/nginx-common/README.Debian
/usr/share/doc/nginx-common/changelog.Debian.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/nginx-common
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/nginx.service
