# rootfs/var/lib/dpkg/info/openbsd-inetd.list
/.
/etc
/etc/init.d
/etc/init.d/openbsd-inetd
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/inetd.service
/usr
/usr/share
/usr/share/man
/usr/share/man/man5
/usr/share/man/man8
/usr/share/man/man8/inetd.8.gz
/usr/share/doc
/usr/share/doc/openbsd-inetd
/usr/share/doc/openbsd-inetd/changelog.Debian.gz
/usr/share/doc/openbsd-inetd/changelog.gz
/usr/share/doc/openbsd-inetd/copyright
/usr/sbin
/usr/sbin/inetd
/lib/systemd/system/openbsd-inetd.service
/usr/share/man/man5/inetd.conf.5.gz
