# rootfs/var/lib/dpkg/info/postfix.list
/.
/etc
/etc/init.d
/etc/init.d/postfix
/etc/ppp
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/postfix
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/postfix
/etc/network
/etc/network/if-up.d
/etc/network/if-up.d/postfix
/etc/network/if-down.d
/etc/network/if-down.d/postfix
/etc/postfix
/etc/postfix/sasl
/etc/postfix/postfix-script
/etc/postfix/postfix-files
/etc/postfix/post-install
/etc/resolvconf
/etc/resolvconf/update-libc.d
/etc/resolvconf/update-libc.d/postfix
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/postfix
/etc/insserv.conf.d
/etc/insserv.conf.d/postfix
/etc/rsyslog.d
/etc/rsyslog.d/postfix.conf
/usr
/usr/share
/usr/share/man
/usr/share/man/man1
/usr/share/man/man1/mailq.1.gz
/usr/share/man/man1/qmqp-sink.1.gz
/usr/share/man/man1/postmulti.1.gz
/usr/share/man/man1/qmqp-source.1.gz
/usr/share/man/man1/postdrop.1.gz
/usr/share/man/man1/posttls-finger.1.gz
/usr/share/man/man1/postqueue.1.gz
/usr/share/man/man1/postfix.1.gz
/usr/share/man/man1/postcat.1.gz
/usr/share/man/man1/newaliases.1.gz
/usr/share/man/man1/postkick.1.gz
/usr/share/man/man1/postmap.1.gz
/usr/share/man/man1/smtp-sink.1.gz
/usr/share/man/man1/sendmail.1.gz
/usr/share/man/man1/postalias.1.gz
/usr/share/man/man1/postconf.1.gz
/usr/share/man/man1/qshape.1.gz
/usr/share/man/man1/postlog.1.gz
/usr/share/man/man1/postsuper.1.gz
/usr/share/man/man1/postlock.1.gz
/usr/share/man/man1/smtp-source.1.gz
/usr/share/man/man8
/usr/share/man/man8/pipe.8postfix.gz
/usr/share/man/man8/smtpd.8postfix.gz
/usr/share/man/man8/qmgr.8postfix.gz
/usr/share/man/man8/scache.8postfix.gz
/usr/share/man/man8/tlsproxy.8postfix.gz
/usr/share/man/man8/discard.8postfix.gz
/usr/share/man/man8/anvil.8postfix.gz
/usr/share/man/man8/trivial-rewrite.8postfix.gz
/usr/share/man/man8/qmqpd.8postfix.gz
/usr/share/man/man8/spawn.8postfix.gz
/usr/share/man/man8/dnsblog.8postfix.gz
/usr/share/man/man8/postfix-add-policy.8.gz
/usr/share/man/man8/showq.8postfix.gz
/usr/share/man/man8/oqmgr.8postfix.gz
/usr/share/man/man8/tlsmgr.8postfix.gz
/usr/share/man/man8/verify.8postfix.gz
/usr/share/man/man8/flush.8postfix.gz
/usr/share/man/man8/pickup.8postfix.gz
/usr/share/man/man8/virtual.8postfix.gz
/usr/share/man/man8/postscreen.8postfix.gz
/usr/share/man/man8/cleanup.8postfix.gz
/usr/share/man/man8/postfix-add-filter.8.gz
/usr/share/man/man8/local.8postfix.gz
/usr/share/man/man8/proxymap.8postfix.gz
/usr/share/man/man8/smtp.8postfix.gz
/usr/share/man/man8/master.8postfix.gz
/usr/share/man/man8/rmail.8.gz
/usr/share/man/man8/lmtp.8postfix.gz
/usr/share/man/man8/error.8postfix.gz
/usr/share/man/man8/bounce.8postfix.gz
/usr/share/man/man5
/usr/share/man/man5/pgsql_table.5.gz
/usr/share/man/man5/relocated.5.gz
/usr/share/man/man5/postfix-wrapper.5.gz
/usr/share/man/man5/access.5.gz
/usr/share/man/man5/aliases.5.gz
/usr/share/man/man5/bounce.5.gz
/usr/share/man/man5/virtual.5.gz
/usr/share/man/man5/generic.5.gz
/usr/share/man/man5/nisplus_table.5.gz
/usr/share/man/man5/canonical.5.gz
/usr/share/man/man5/header_checks.5.gz
/usr/share/man/man5/pcre_table.5.gz
/usr/share/man/man5/memcache_table.5.gz
/usr/share/man/man5/postconf.5.gz
/usr/share/man/man5/body_checks.5.gz
/usr/share/man/man5/regexp_table.5.gz
/usr/share/man/man5/ldap_table.5.gz
/usr/share/man/man5/sqlite_table.5.gz
/usr/share/man/man5/mysql_table.5.gz
/usr/share/man/man5/tcp_table.5.gz
/usr/share/man/man5/lmdb_table.5.gz
/usr/share/man/man5/transport.5.gz
/usr/share/man/man5/cidr_table.5.gz
/usr/share/man/man5/socketmap_table.5.gz
/usr/share/man/man5/master.5.gz
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/postfix
/usr/share/doc
/usr/share/doc/postfix
/usr/share/doc/postfix/changelog.gz
/usr/share/doc/postfix/changelog.Debian.gz
/usr/share/doc/postfix/README.Debian
/usr/share/doc/postfix/NEWS.Debian.gz
/usr/share/doc/postfix/copyright
/usr/share/postfix
/usr/share/postfix/main.cf.debian
/usr/share/postfix/main.cf.tls
/usr/share/postfix/postinst.functions
/usr/share/postfix/main.cf.dist
/usr/share/postfix/master.cf.dist
/usr/lib
/usr/lib/libpostfix-global.so.1.0.1
/usr/lib/libpostfix-dns.so.1.0.1
/usr/lib/postfix
/usr/lib/postfix/tlsproxy
/usr/lib/postfix/spawn
/usr/lib/postfix/smtpd
/usr/lib/postfix/flush
/usr/lib/postfix/scache
/usr/lib/postfix/virtual
/usr/lib/postfix/tlsmgr
/usr/lib/postfix/pickup
/usr/lib/postfix/master
/usr/lib/postfix/smtp
/usr/lib/postfix/showq
/usr/lib/postfix/local
/usr/lib/postfix/cleanup
/usr/lib/postfix/postscreen
/usr/lib/postfix/discard
/usr/lib/postfix/master.cf
/usr/lib/postfix/proxymap
/usr/lib/postfix/postfix-script
/usr/lib/postfix/bounce
/usr/lib/postfix/qmgr
/usr/lib/postfix/dnsblog
/usr/lib/postfix/error
/usr/lib/postfix/fsstone
/usr/lib/postfix/verify
/usr/lib/postfix/postfix-files
/usr/lib/postfix/postmulti-script
/usr/lib/postfix/post-install
/usr/lib/postfix/postfix-wrapper
/usr/lib/postfix/pipe
/usr/lib/postfix/anvil
/usr/lib/postfix/qmqpd
/usr/lib/postfix/dict_tcp.so
/usr/lib/postfix/dict_sqlite.so
/usr/lib/postfix/main.cf
/usr/lib/postfix/trivial-rewrite
/usr/lib/postfix/oqmgr
/usr/lib/postfix_groups.pl
/usr/lib/libpostfix-master.so.1.0.1
/usr/lib/libpostfix-tls.so.1.0.1
/usr/lib/libpostfix-util.so.1.0.1
/usr/sbin
/usr/sbin/postalias
/usr/sbin/postmulti
/usr/sbin/qmqp-source
/usr/sbin/postqueue
/usr/sbin/smtp-sink
/usr/sbin/postfix-add-filter
/usr/sbin/postsuper
/usr/sbin/sendmail
/usr/sbin/qshape
/usr/sbin/postconf
/usr/sbin/posttls-finger
/usr/sbin/rmail
/usr/sbin/smtp-source
/usr/sbin/postdrop
/usr/sbin/postlock
/usr/sbin/postfix
/usr/sbin/postcat
/usr/sbin/qmqp-sink
/usr/sbin/postmap
/usr/sbin/postlog
/usr/sbin/postkick
/usr/sbin/postfix-add-policy
/usr/bin
/var
/var/spool
/var/spool/postfix
/var/spool/postfix/lib
/var/spool/postfix/etc
/var/spool/postfix/usr
/var/spool/postfix/usr/lib
/var/spool/postfix/usr/lib/zoneinfo
/var/spool/postfix/usr/lib/sasl2
/var/spool/postfix/dev
/var/lib
/var/lib/postfix
/var/log
/usr/share/man/man8/trace.8postfix.gz
/usr/share/man/man8/defer.8postfix.gz
/usr/lib/libpostfix-util.so.1
/usr/lib/libpostfix-dns.so.1
/usr/lib/sendmail
/usr/lib/libpostfix-master.so.1
/usr/lib/libpostfix-global.so.1
/usr/lib/postfix/lmtp
/usr/lib/libpostfix-tls.so.1
/usr/bin/mailq
/usr/bin/rmail
/usr/bin/newaliases
