# rootfs/var/lib/dpkg/info/pure-ftpd.list
/.
/etc
/etc/ftpusers
/etc/init.d
/etc/init.d/pure-ftpd
/usr
/usr/sbin
/usr/sbin/pure-ftpwho
/usr/sbin/pure-uploadscript
/usr/sbin/pure-ftpd
/usr/sbin/pure-mrtginfo
/usr/sbin/pure-ftpd-virtualchroot
/usr/sbin/pure-quotacheck
/usr/sbin/pure-authd
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/pure-ftpd
/usr/share/doc
/usr/bin
/usr/bin/pure-pw
/usr/bin/pure-pwconvert
/usr/bin/pure-statsdecode
/usr/share/man/man8/pure-ftpd-virtualchroot.8.gz
/usr/share/doc/pure-ftpd
