# rootfs/etc/cron.d/csf_update
SHELL=/bin/sh
32 23 * * * root /usr/sbin/csf -u
